# profile.d/history.sh
export HISTTIMEFORMAT='%y/%m/%d %H:%M:%S '
export HISTSIZE=10000
